# autoupdate-unit/usr/lib/systemd/system/autoupdate-host.timer
[Unit]
Description=Automatic daily host upgrades
[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d
[Install]
WantedBy=timers.target
